# supervisord.conf
[unix_http_server]
file = /tmp/supervisor.sock

[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = /var/run/supervisord.pid
nodaemon = true
minfds = 1024
minprocs = 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:sample]
command = python sample.py
autostart = false
stdout_logfile = /tmp/svsout.log
directory = /mnt/app/sampleapp
stopsignal = QUIT
startretries = 0

[program:nginx]
command = nginx
autostart = false
stdout_logfile = /tmp/nginx.out
directory = /mnt/app
stopsignal = QUIT
startretries = 0

[program:uwsgi]
command = sudo -u www-data uwsgi -s /tmp/uwsgi.sock -w sample:app
autostart = false
stdout_logfile = /tmp/uwsgi.log
directory = /mnt/app/sampleapp
stopsignal = QUIT
startretries = 0

